# jasonernst.com/conf/supervisor/geth.conf
[program:geth]
command=geth
user=ethminer
stdout_logfile=/var/log/supervisor/geth-stdout.log
stderr_logfile=/var/log/supervisor/geth-stderr.log
