# jasonernst.com/conf/supervisor/golem.conf
[program:golem]
directory=/home/jason/golem
command=bash golemstart.sh
user=jason
stdout_logfile=/var/log/supervisor/golem-stdout.log
stderr_logfile=/var/log/supervisor/golem-stderr.log
