# jasonernst.com/conf/supervisor/ipfs.conf
[program:ipfs]
command=/usr/local/bin/ipfs daemon
user=jason
stdout_logfile=/var/log/supervisor/ipfs-stdout.log
stderr_logfile=/var/log/supervisor/ipfs-stderr.log
