# ipfs/ipfs.service
[Unit]
Description=Ipfs

[Service]
ExecStart="/opt/IPFS Desktop/ipfs-desktop" --no-sandbox %U

[Install]
WantedBy=multi-user.target
